# secmon-linux-security/systemd/secmon-api.service
[Unit]
Description=SecMon API Service
After=network.target

[Service]
Type=simple
User=secmon
Group=secmon
WorkingDirectory=/opt/secmon
EnvironmentFile=/etc/secmon/secmon.env
ExecStart=/opt/secmon/.venv/bin/uvicorn backend.app:app --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=300
StartLimitBurst=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/secmon /var/log/secmon

[Install]
WantedBy=multi-user.target
